# systemd/hermes-backup-prune.service
[Unit]
Description=Monthly restic forget+prune (apply retention policy)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=admin
Group=admin
EnvironmentFile=/home/admin/code/hermes-config/.env.backup
ExecStart=/usr/local/bin/restic --option b2.connections=4 forget --keep-daily 14 --keep-weekly 8 --keep-monthly 12 --keep-yearly 3 --prune
TimeoutStartSec=12h

[Install]
WantedBy=multi-user.target
